# autostart/lxqt-qlipper-autostart.desktop.in
[Desktop Entry]
Type=Application
Name=Qlipper
TryExec=qlipper
Exec=qlipper
OnlyShowIn=LXQt;
X-LXQt-Need-Tray=true

#TRANSLATIONS_DIR=translations
